# systemd/agent-co-auth-admin.service
[Unit]
Description=agent-co Auth Admin Service
After=network.target

[Service]
Type=simple
WorkingDirectory=__APP_DIR__
EnvironmentFile=-/etc/agent-co/agent-co-auth-admin.env
ExecStart=/usr/bin/node __APP_DIR__/dist/auth-admin-server.js
Restart=always
RestartSec=2
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=__APP_DIR__/data

[Install]
WantedBy=multi-user.target
