# systemd/agent-co.service
[Unit]
<<<<<<<< HEAD:systemd/agent-co.service
Description=agent-co (chat + auth-admin)
After=network.target redis-server.service
Requires=redis-server.service
========
Description=agent-co Chat Service
After=network.target redis-server.service agent-co-auth-admin.service
Requires=redis-server.service agent-co-auth-admin.service
>>>>>>>> origin/main:systemd/agent-co-chat.service

[Service]
Type=simple
WorkingDirectory=__APP_DIR__
EnvironmentFile=-/etc/agent-co/agent-co.env
ExecStartPre=/bin/sh -c 'redis-cli -u "redis://127.0.0.1:6379" ping | grep -q PONG'
ExecStart=/bin/bash __APP_DIR__/scripts/start-services.sh
Restart=always
RestartSec=2
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=__APP_DIR__/data
ReadOnlyPaths=__APP_DIR__/public

[Install]
WantedBy=multi-user.target
