# systemd/bot-room-chat.service
[Unit]
Description=Bot Room Chat Service
After=network.target bot-room-auth-admin.service
Requires=bot-room-auth-admin.service

[Service]
Type=simple
WorkingDirectory=/opt/bot-room
EnvironmentFile=-/etc/bot-room/bot-room-chat.env
ExecStart=/usr/bin/node /opt/bot-room/dist/server.js
Restart=always
RestartSec=2
User=botroom
Group=botroom
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadOnlyPaths=/opt/bot-room/public

[Install]
WantedBy=multi-user.target
